# uart/bin/99-usb-logger.rules
ACTION=="add", SUBSYSTEM=="tty", RUN+="/usr/local/bin/p1.sh %k"
