# uart/bin/99-usb.rules
ACTION=="add", SUBSYSTEM=="usb", RUN+="/bin/bash /usr/local/bin/p2.sh"
