# uart/bin/99-usb-logger.rules
ACTION=="add", SUBSYSTEM=="tty", RUN+="/bin/sh -c '/usr/local/bin/p1.sh %k | at now'"